# deploy/ai.nimbus.weixin-deepseek.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ai.nimbus.weixin-deepseek.plist — launchd agent for the nimbus WeChat/DeepSeek instance.

  This is a TEMPLATE.  Fill in the placeholder values (marked with <FILL_ME>)
  before installing.  Do NOT commit with real secrets.

  INSTALL (manual — after filling in secrets):
    cp ~/nimbus-os/nimbus/deploy/ai.nimbus.weixin-deepseek.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/ai.nimbus.weixin-deepseek.plist

  UNLOAD:
    launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/ai.nimbus.weixin-deepseek.plist

  ⚠️  Isolation guarantees:
    1. PROVIDER=deepseek — activates the DeepSeek backend; Cici's plist has no
       PROVIDER key, so it always runs on the Claude path.
    2. ANTHROPIC_BASE_URL + ANTHROPIC_API_KEY are scoped to THIS plist only.
       Never put them in ~/.zshrc, ~/.env, or Cici's plist.
    3. NIMBUS_DB_PATH points to a separate SQLite file — no session sharing
       with Cici's data/state.db.
    4. NIMBUS_OUTBOX_DIR points to a separate outbox directory so chart files
       from the two instances never mix.
    5. wechat-io POSTs OpenAI-compatible requests to the WeixinInboundChannel on
       WEIXIN_INBOUND_PORT (8770).  This is the /v1/chat/completions endpoint —
       NOT the local ApiChannel (/chat) which is disabled here.  8770 avoids the
       8788 default used by the Phase-2 iLink two-way channel.
       After installing, update ~/.wechat-io/config.yaml:
         agent:
           openai_base_url: "http://127.0.0.1:8770/v1"   # → POST /v1/chat/completions
           openai_api_key: "<same value as WEIXIN_INBOUND_TOKEN below>"
           model: "deepseek-chat"
       And send X-Hermes-Session-Id: wechat-{chat_id} for per-group session isolation.

  ⚠️  Single-consumer rule:
    Confirm that Hermes is stopped (launchctl bootout gui/$(id -u) com.hermes.daemon)
    before enabling this service — both consume wechat-io's webhook and running both
    causes double replies and session conflicts.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.nimbus.weixin-deepseek</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/x/nimbus-os/nimbus/scripts/nimbus-daemon-weixin.sh</string>
  </array>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/x/nimbus-os/nimbus/logs/weixin-deepseek.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/x/nimbus-os/nimbus/logs/weixin-deepseek.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- miniforge first (futu-capable python3), then ~/.bun/bin (bun/tmux). -->
    <key>PATH</key>
    <string>/Users/x/miniforge3/bin:/Users/x/.bun/bin:/usr/local/bin:/usr/bin:/bin</string>

    <!-- ── Provider selection ── -->
    <!-- Activates the DeepSeek code path throughout nimbus. -->
    <key>PROVIDER</key>
    <string>deepseek</string>

    <!-- ── DeepSeek API credentials ── -->
    <!-- Point the Agent SDK at DeepSeek's Anthropic-compatible endpoint. -->
    <key>ANTHROPIC_BASE_URL</key>
    <string>https://api.deepseek.com/anthropic</string>
    <!-- Replace with your real DeepSeek API key (sk-...). -->
    <key>ANTHROPIC_API_KEY</key>
    <string><FILL_ME:sk-deepseek-your-key-here></string>

    <!-- ── Instance isolation ── -->
    <!-- Separate SQLite DB — must NOT be the same file as Cici's data/state.db. -->
    <key>NIMBUS_DB_PATH</key>
    <string>/Users/x/nimbus-os/nimbus/data/weixin-deepseek.db</string>
    <!-- Separate outbox so chart files from the two instances don't mix. -->
    <key>NIMBUS_OUTBOX_DIR</key>
    <string>/Users/x/nimbus-os/nimbus/data/outbox-weixin</string>

    <!-- ── WeChat / wechat-io connector (OpenAI-compatible inbound) ── -->
    <!-- Turn on WeixinInboundChannel; it serves POST /v1/chat/completions. -->
    <key>WEIXIN_INBOUND</key>
    <string>1</string>
    <!-- Port wechat-io POSTs to.  Must match ~/.wechat-io/config.yaml
         (agent.openai_base_url = http://127.0.0.1:8770/v1).  8770 avoids the
         8788 default used by the Phase-2 iLink two-way channel. -->
    <key>WEIXIN_INBOUND_PORT</key>
    <string>8770</string>
    <!-- Shared bearer secret for /v1/chat/completions.  Put the SAME value in
         wechat-io's config as openai_api_key. -->
    <key>WEIXIN_INBOUND_TOKEN</key>
    <string><FILL_ME:any-random-secret></string>
    <!-- Disable Discord — this instance only serves wechat-io (has no bot token). -->
    <key>NIMBUS_DISCORD_ENABLED</key>
    <string>0</string>
    <!-- Disable the local ApiChannel (/chat); wechat-io uses the OpenAI endpoint. -->
    <key>NIMBUS_API_ENABLED</key>
    <string>0</string>
    <!-- Disable the iLink two-way WeChat channel (Phase 2 — not used here). -->
    <key>WEIXIN_TWOWAY</key>
    <string>0</string>
    <!-- Disable WeChat mirror pushes (this instance IS the WeChat instance). -->
    <key>WEIXIN_MIRROR</key>
    <string>0</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/x/nimbus-os/nimbus</string>
</dict>
</plist>
