# deploy/com.local.masters-study.review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 22:30 CST — 跨市场复盘 → 按星期轮换一位大师 -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.masters-study.review</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/x/nimbus-os/nimbus/scripts/masters-study-digest.sh</string>
    <string>review</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/x/nimbus-os/nimbus/logs/masters-study.review.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/x/nimbus-os/nimbus/logs/masters-study.review.err.log</string>
</dict>
</plist>
