# skills/square-bot/bot/launchd/com.local.binance-square-bot.generate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.binance-square-bot.generate</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/x/.claude/skills/square-bot/bot/scripts/run_generation_cycle.sh</string>
    <string>all</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>55</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>35</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>5</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>16</integer>
      <key>Minute</key>
      <integer>35</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>19</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>22</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/x/.claude/skills/square-bot/bot</string>

  <key>StandardOutPath</key>
  <string>/Users/x/.claude/skills/square-bot/bot/logs/launchd-generate.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/x/.claude/skills/square-bot/bot/logs/launchd-generate.err.log</string>
</dict>
</plist>
